# design/remote_access/.config/autostart/shades.desktop
[Desktop Entry]
Name=shades
Type=Application
Categories=Application
Exec=/home/pi/shades.sh
Terminal=true
StartupNotify=false
X-KeepTerminal=true
